# adb/adb-udev.rules
SUBSYSTEM=="udc",ACTION=="change",DRIVER=="configfs-gadget",RUN+="/usr/bin/adb_config %E{DEVPATH}"
